# etc/network/if-pre-up.d/start-hostapd
#!/bin/sh
if [ "$IFACE" = wlan0 ]; then
	killall -9 hostapd
	sleep 1
	/usr/local/bin/hostapd /etc/hostapd/hostapd.conf &> /dev/null &
fi
exit 0
